# systemd/claude-code-schedule.user.service
[Unit]
Description=Claude Code Schedule (user) - Automated Claude Code execution
After=network.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/root1/githubTools/claude-code-schedule
# Ensure the Claude CLI is found under systemd user session
Environment=CLAUDE_BIN=/home/root1/.nvm/versions/node/v22.18.0/bin/claude
Environment=PATH=/home/root1/.nvm/versions/node/v22.18.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
ExecStart=/home/root1/githubTools/claude-code-schedule/target/release/ccschedule --loop-mode --ping-mode --log-dir /home/root1/githubTools/claude-code-schedule/logs
Restart=always
RestartSec=5

[Install]
WantedBy=default.target

