# launchd/com.claude-code-schedule.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-code.schedule.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/REPLACE_WITH_USERNAME/githubTools/claude-code-schedule/target/release/ccschedule</string>
        <string>--loop-mode</string>
        <string>--ping-mode</string>
        <string>--log-dir</string>
        <string>/Users/REPLACE_WITH_USERNAME/Library/Logs/claude-code-schedule</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/REPLACE_WITH_USERNAME/githubTools/claude-code-schedule</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CLAUDE_BIN</key>
        <string>/opt/homebrew/bin/claude</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/REPLACE_WITH_USERNAME/Library/Logs/claude-code-schedule/agent.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/REPLACE_WITH_USERNAME/Library/Logs/claude-code-schedule/agent.err.log</string>
</dict>
</plist>

